# scripts/launchd/org.solutionsmith.its.picklist-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Picklist sync — hourly (StartInterval 3600; see the operator-decision note below).

  Reads enabled mappings from Picklist_Sync_Config; for each, syncs
  source-master-DB unique values into the target sheet's PICKLIST
  options. Idempotent via last_run_hash on each mapping row.

  Install (renders the __ITS_HOME__ placeholder via install.sh — do NOT `cp`
  this file directly; a raw copy would install an unresolved placeholder path):
    cd scripts/launchd && ./install.sh load org.solutionsmith.its.picklist-sync

  This file uses the __ITS_HOME__ template placeholder like every other daemon
  plist (install.sh substitutes $HOME/its at render time). It was the last plist
  carrying literal /Users/sethsmith paths; templatized 2026-07-12 so the fleet is
  uniform and portable to a new host/username (host-migration robustness).

  See docs/picklist_sync.md for the full operations runbook.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.solutionsmith.its.picklist-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>__ITS_HOME__/.venv/bin/python</string>
        <string>__ITS_HOME__/scripts/run_picklist_sync.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__ITS_HOME__</string>

    <key>StandardOutPath</key>
    <string>__ITS_HOME__/logs/launchd/run_picklist_sync.out.log</string>

    <key>StandardErrorPath</key>
    <string>__ITS_HOME__/logs/launchd/run_picklist_sync.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <!-- Recurring every hour. Operator decision 2026-05-21: hourly
         matches construction-workflow tempo (vendors added today don't
         need to appear in dropdowns within 15 minutes), cuts run count
         96/day → 24/day, and the idempotency short-circuit means most
         hourly runs are still cheap reads. -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- A2 (single-host resilience): RunAtLoad=true so this interval poll daemon
         RESUMES immediately when launchd loads the agent after a reboot/OS-update
         (eval host-no-runatload), then StartInterval re-fires on cadence. KeepAlive
         is intentionally absent: ITS daemons are one-shot-per-fire, so KeepAlive=true
         would restart on exit and destroy the interval cadence. -->
    <key>RunAtLoad</key>
    <true/>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
