# scripts/launchd/org.solutionsmith.its.rfq-send.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  RFQ SEND poll daemon (ADR-0004 R3, RFQ/estimate sub-lane of po_materials). The send half
  of the External Send Gate two-process model for the VENDOR audience: discovers approved
  RFQ_Pending_Review rows (Send Now / Approve for Scheduled Send), runs the F22 approval-
  attestation gate against the ITS — Purchase Orders workspace (§46 — the SAME procurement
  approver set as POs), and dispatches each to po_materials.rfq_send.send_one_row (→ the
  shared weekly_send engine → Graph send from procurement@). The RFQ email carries TWO
  attachments: the PRICE-FREE RFQ PDF + the vendor's fillable .xlsx quote form (the R3
  sequence-attachment engine seam). AI-FREE — anthropic / anthropic_client AND ollama_client
  AST-forbidden in tests/test_capability_gating.py SEND_SCRIPTS (send scripts are
  local-AI-free too, ADR-0004 decision 12); the generation side (rfq-poll) is a SEPARATE
  process with zero send capability. Recipient resolved at send time from ITS_Vendors by
  Vendor Key. Tagged po_materials_rfq so cross-lane dispatch (po-send / subcontract-send) is
  structurally impossible.

  SHIPS DARK. Gated OFF by po_materials.rfq_send.polling_enabled
  (scripts/migrations/seed_rfq_send_config.py seeds it false). Go-live is a FIXED
  high-capability-class External-Send-Gate operator action (Seth): flip
  po_materials.rfq_send.polling_enabled true AND load this plist together.

  Install (go-live only):
    scripts/launchd/install.sh load org.solutionsmith.its.rfq-send

  StartInterval is PARAMETERIZED (install.sh __POLL_INTERVAL_SECONDS__): resolution is
  [interval] arg > ITS_Config po_materials.rfq_send.poll_interval_seconds > default 900s
  (15 min, matching weekly-send / progress-send / po-send / subcontract-send — an approval
  poller, not a fast puller). The resolved value is BAKED into the installed plist, so a
  later change needs a re-load.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.solutionsmith.its.rfq-send</string>

    <key>ProgramArguments</key>
    <array>
        <string>__ITS_HOME__/.venv/bin/python</string>
        <string>-m</string>
        <string>po_materials.rfq_send_poll</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__ITS_HOME__</string>

    <key>StandardOutPath</key>
    <string>__ITS_HOME__/logs/launchd/rfq_send.out.log</string>

    <key>StandardErrorPath</key>
    <string>__ITS_HOME__/logs/launchd/rfq_send.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>StartInterval</key>
    <integer>__POLL_INTERVAL_SECONDS__</integer>

    <!-- A2 (single-host resilience): RunAtLoad=true so this interval daemon RESUMES
         immediately when launchd loads the agent after a reboot/OS-update, then
         StartInterval re-fires on cadence. KeepAlive is intentionally absent: ITS daemons
         are one-shot-per-fire, so KeepAlive=true would restart on exit and destroy the
         interval cadence. -->
    <key>RunAtLoad</key>
    <true/>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
