# scripts/launchd/org.solutionsmith.its.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ITS daily watchdog. Per Operational Standards v8 §2.
  Runs scripts/watchdog.py at 7:00 AM local time. launchd's
  StartCalendarInterval catches up on wake, so if the laptop was
  asleep at 7:00 AM the job runs as soon as it wakes.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.solutionsmith.its.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>__ITS_HOME__/.venv/bin/python</string>
        <string>__ITS_HOME__/scripts/watchdog.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__ITS_HOME__</string>

    <key>StandardOutPath</key>
    <string>__ITS_HOME__/logs/launchd/watchdog.out.log</string>

    <key>StandardErrorPath</key>
    <string>__ITS_HOME__/logs/launchd/watchdog.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- A2 (single-host resilience): RunAtLoad stays FALSE for a StartCalendarInterval
         daemon. launchd fires a missed calendar job on wake; RunAtLoad=true would
         MIS-fire it on every login (e.g. run a Friday compile on a Tuesday boot).
         Reboot-recovery for the weekly compile is the watchdog Check-I catch-up. -->
    <key>RunAtLoad</key>
    <false/>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
