# scripts/launchd/org.solutionsmith.its.portal-poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Safety Portal pull-model polling daemon. Drains the Cloudflare Worker's D1
  submission queue: GET /api/internal/pending -> per-row canonical HMAC verify
  (downgrade defense) -> safety_reports.intake.process_portal_submission ->
  POST /api/internal/mark-filed. Zero external-send capability (enrolled in
  tests/test_capability_gating.py GATED_SCRIPTS; egress only via the
  F02-allowlisted shared.portal_client). Fail-closed on missing creds.

  Install:
    scripts/launchd/install.sh load org.solutionsmith.its.portal-poll

  The installer reads safety_reports.portal_poll.poll_interval_seconds from
  ITS_Config and substitutes both __ITS_HOME__ and __POLL_INTERVAL_SECONDS__
  before copying to ~/Library/LaunchAgents/. Default cadence 60s. Changes to the
  interval row take effect at the next install run; the running daemon does not
  hot-reload (launchd holds the interval).

  See safety_reports/README.md §43 for the operator runbook.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.solutionsmith.its.portal-poll</string>

    <key>ProgramArguments</key>
    <array>
        <string>__ITS_HOME__/.venv/bin/python</string>
        <string>-m</string>
        <string>safety_reports.portal_poll</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__ITS_HOME__</string>

    <key>StandardOutPath</key>
    <string>__ITS_HOME__/logs/launchd/portal_poll.out.log</string>

    <key>StandardErrorPath</key>
    <string>__ITS_HOME__/logs/launchd/portal_poll.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>StartInterval</key>
    <integer>__POLL_INTERVAL_SECONDS__</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
